# etc/rc.local
#!/bin/sh -e

/opt/ghpi/scripts/sensors.py >/dev/null 2>&1 &
#/home/pi/ghpi/remotemon.py >/dev/null 2>&1 &

exit 0
