# etc/rc.local
#!/bin/sh -e

/home/pi/ghpi/tempmon.py >/dev/null 2>&1 &
/home/pi/ghpi/remotemon.py >/dev/null 2>&1 &

exit 0
